# packaging/deb/postinst
#!/bin/bash
set -e

echo "Setting up mcp-bash-server..."

# Create user if it doesn't exist
if ! id -u mcp >/dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin mcp
    echo "Created user: mcp"
fi

# Handle config.toml: never overwrite existing config
CONFIG_DIR=/etc/mcp-bash-server
CONFIG_FILE=$CONFIG_DIR/config.toml
EXAMPLE_FILE=$CONFIG_DIR/config.example.toml

if [ -f "$CONFIG_FILE" ]; then
    echo "Existing config.toml found — preserving your settings."
    echo ""
    echo "=== Configuration changes in this version ==="
    if [ -f "$EXAMPLE_FILE" ]; then
        diff "$CONFIG_FILE" "$EXAMPLE_FILE" 2>/dev/null || true
    fi
    echo "=== Review config.example.toml for new options ==="
    echo ""
else
    # First install: copy example as initial config
    cp "$EXAMPLE_FILE" "$CONFIG_FILE"

    # Generate random API key
    API_KEY=$(openssl rand -hex 24 2>/dev/null || head -c 32 /dev/urandom | xxd -p | head -c 48)
    if [ -n "$API_KEY" ]; then
        sed -i "s/^# api_key = .*/api_key = \"$API_KEY\"/" "$CONFIG_FILE" 2>/dev/null || \
        sed -i "s/^api_key = .*/api_key = \"$API_KEY\"/" "$CONFIG_FILE" 2>/dev/null || true
        echo ""
        echo "============================================"
        echo "  Generated API key:"
        echo "  $API_KEY"
        echo "============================================"
        echo "  Save this key — it's shown only once!"
        echo "  Use it in MCP client config as Bearer token."
        echo "============================================"
        echo ""
    fi

    # Default bind to localhost for security
    sed -i 's/^host = "0.0.0.0"/host = "127.0.0.1"/' "$CONFIG_FILE" 2>/dev/null || true
    echo "Default bind: 127.0.0.1 (localhost only)."
    echo "Change to 0.0.0.0 in config.toml for remote access."
    echo ""
fi

chown mcp:mcp "$CONFIG_FILE"
chmod 640 "$CONFIG_FILE"

# Create process data directory
PROCESS_DIR=/var/lib/mcp-bash-server
mkdir -p "$PROCESS_DIR"/output
chown -R mcp:mcp "$PROCESS_DIR"
chmod 700 "$PROCESS_DIR"
chmod 700 "$PROCESS_DIR"/output

# Ask about sudo for mcp user (only on first install, interactive)
OVERRIDE_DIR=/etc/systemd/system/mcp-bash-server.service.d
OVERRIDE_FILE="$OVERRIDE_DIR/sudo-override.conf"
SUDOERS_FILE=/etc/sudoers.d/mcp

setup_sudo() {
    echo ""
    echo "=== Sudo Access Configuration ==="
    echo "Granting sudo to the mcp user allows AI agents to run commands"
    echo "with root privileges (apt, systemctl, editing system files)."
    echo ""
    echo "Without sudo, agents are limited to the mcp user's permissions."
    echo "This also affects systemd security sandboxing:"
    echo "  - With sudo:    sandbox is disabled (sudo won't work in sandbox)"
    echo "  - Without sudo: sandbox is enabled (recommended for security)"
    echo ""

    if [ ! -t 0 ]; then
        echo "Non-interactive install detected."
        echo "To enable sudo later, run:"
        echo "  sudo /usr/sbin/mcp-bash-server-configure sudo"
        echo ""
        return 1
    fi

    while true; do
        read -p "Grant sudo (root) access to the mcp user? [y/N] " yn
        case "$yn" in
            [Yy]|[Yy][Ee][Ss]) return 0 ;;
            [Nn]|[Nn][Oo]|"") return 1 ;;
            *) echo "Please answer yes or no." ;;
        esac
    done
}

if setup_sudo; then
    # Grant sudo
    echo "mcp ALL=(ALL) NOPASSWD: ALL" > "$SUDOERS_FILE"
    chmod 440 "$SUDOERS_FILE"
    chown root:root "$SUDOERS_FILE"

    # Disable systemd sandboxing via override (not modifying the packaged unit)
    mkdir -p "$OVERRIDE_DIR"
    cat > "$OVERRIDE_FILE" <<'EOF'
# Sudo access requires disabling some systemd security restrictions
[Service]
NoNewPrivileges=
ProtectSystem=full
ProtectHome=
ReadWritePaths=/tmp /var/lib/mcp-bash-server
PrivateTmp=
ProtectKernelTunables=
ProtectKernelModules=
ProtectControlGroups=
RestrictSUIDSGID=
RestrictRealtime=
RestrictNamespaces=
LockPersonality=
MemoryDenyWriteExecute=
SystemCallFilter=
SystemCallErrorNumber=
EOF
    echo "Sudo access granted. Systemd sandboxing disabled via override."
else
    # Remove sudo if it was previously granted
    rm -f "$SUDOERS_FILE" 2>/dev/null || true

    # Remove sandbox override if it existed
    rm -f "$OVERRIDE_FILE" 2>/dev/null || true
    if [ -d "$OVERRIDE_DIR" ] && [ -z "$(ls -A "$OVERRIDE_DIR" 2>/dev/null)" ]; then
        rmdir "$OVERRIDE_DIR" 2>/dev/null || true
    fi

    echo "No sudo access. Systemd sandboxing enabled (recommended)."
fi

# Reload systemd
systemctl daemon-reload

echo ""
echo "mcp-bash-server installed successfully!"
echo ""
echo "Next steps:"
echo "  1. Review config changes: diff /etc/mcp-bash-server/config.toml /etc/mcp-bash-server/config.example.toml"
echo "  2. Edit config: nano /etc/mcp-bash-server/config.toml"
echo "  3. Set API key: export MCP_API_KEY=your-secret-key"
echo "  4. Start the service: systemctl enable --now mcp-bash-server"
echo ""
echo "To reconfigure sudo access later: sudo mcp-bash-server-configure sudo"
echo "Check status: systemctl status mcp-bash-server"
echo "View logs: journalctl -u mcp-bash-server -f"
